# supervisor-app.conf
[program:app-uwsgi]
command = /usr/bin/uwsgi --ini /opt/flask/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx
